# platforms/organelle_m/rootfs/usr/lib/systemd/system/cherrypy.service
[Unit] 
Description=Cherrypy Web Server

[Service] 
Type=simple 
WorkingDirectory=/home/patch/fw_dir/web/server
ExecStart=/home/patch/fw_dir/web/server/run.sh
User=patch
Group=User=patch


[Install] 
WantedBy=multi-user.target
